# anxious-scroll-daemon.service
[Unit]
Description=Anxious Scroll Daemon - Mouse scroll wheel acceleration
Documentation=https://github.com/Snehal-Reddy/cursor-anxious
After=multi-user.target
Wants=multi-user.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/local/bin/anxious-scroll-daemon
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=anxious-scroll-daemon

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/dev/input

# Required capabilities for device access
CapabilityBoundingSet=CAP_SYS_ADMIN
AmbientCapabilities=CAP_SYS_ADMIN

[Install]
WantedBy=multi-user.target
